# scripts/rpm/post-install.sh
#!/bin/bash

npm install -g serve

update-rc.d setter-for-catan defaults
